# tools/telegram_bridge/tg-bot.service
[Unit]
Description=Telegram Agent Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/opencode/tools/telegram_bridge
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
StandardOutput=append:bot.log
StandardError=append:bot.log

[Install]
WantedBy=multi-user.target
