# tgbot/cipher-bot.service
[Unit]
Description=Cipher Telegram Bot
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/opencode_base
ExecStart=/home/admin/opencode_base/.venv/bin/python3 tgbot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
